# deploy/aicrm-archive-sync.service
[Unit]
Description=AI-CRM Next WeCom archive sync
After=network.target openclaw-wecom-postgres.service
Requires=openclaw-wecom-postgres.service

[Service]
Type=oneshot
User=ubuntu
EnvironmentFile=/home/ubuntu/.openclaw-wecom-pg.env
Environment=APP_HOST=127.0.0.1
Environment=APP_PORT=5001
Environment=WECOM_ARCHIVE_SYNC_LIMIT=100
Environment=WECOM_ARCHIVE_SYNC_MAX_PAGES=200
Environment=WECOM_ARCHIVE_SYNC_MODE=direct
Environment=WECOM_ARCHIVE_SYNC_TIMEOUT_SECONDS=900
WorkingDirectory=/home/ubuntu/极简 crm
ExecStart=/bin/bash -lc 'source /home/ubuntu/venvs/openclaw/bin/activate && python -m scripts.run_incremental_archive_sync'
