# deploy/aicrm-reply-monitor-capture.timer
[Unit]
Description=Run AI-CRM Reply Monitor Capture every 3 minutes

[Timer]
OnCalendar=*-*-* *:00/3:00
Persistent=true
Unit=aicrm-reply-monitor-capture.service

[Install]
WantedBy=timers.target
