# deploy/aicrm-wechat-shop-order-sync.timer
[Unit]
Description=Run AI-CRM WeChat Shop order sync every 10 minutes

[Timer]
OnCalendar=*-*-* *:0/10:30
Persistent=true
Unit=aicrm-wechat-shop-order-sync.service

[Install]
WantedBy=timers.target
