# deploy/openclaw-external-contact-sync.timer
[Unit]
Description=Run OpenClaw WeCom external contact incremental sync every 5 minutes

[Timer]
OnCalendar=*:0/5
Persistent=true
Unit=openclaw-external-contact-sync.service

[Install]
WantedBy=timers.target
