# deploy/openclaw-external-push-worker.timer
[Unit]
Description=Run OpenClaw External Push Worker every minute

[Timer]
OnCalendar=*-*-* *:*:20
Persistent=true
Unit=openclaw-external-push-worker.service

[Install]
WantedBy=timers.target
