# deploy/openclaw-external-contact-sync.timer
[Unit]
Description=Run OpenClaw WeCom external contact incremental sync every 15 minutes

[Timer]
OnCalendar=*-*-* *:05,20,35,50:00
Persistent=true
Unit=openclaw-external-contact-sync.service

[Install]
WantedBy=timers.target
