# deploy/openclaw-ai-audience-scheduler.service
[Unit]
Description=OpenClaw AI Audience Scheduler - emit refresh ticks and dispatch AI audience internal consumers
After=network.target openclaw-wecom-postgres.service
Requires=openclaw-wecom-postgres.service

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/home/ubuntu/极简 crm
EnvironmentFile=/home/ubuntu/.openclaw-wecom-pg.env
Environment=AICRM_INTERNAL_EVENTS_ENABLED=1
Environment=AICRM_INTERNAL_EVENTS_SHADOW_ONLY=0
Environment=AICRM_INTERNAL_EVENTS_AUTO_EXECUTE=1
Environment=AICRM_INTERNAL_EVENTS_AUTO_EXECUTE_MAX_BATCH_SIZE=20
Environment=AICRM_INTERNAL_EVENTS_ALLOWED_EVENT_TYPES=ai_audience.source.changed,channel_entry.entered,questionnaire.submitted,external_form.submitted,payment.succeeded,transaction.paid,payment_succeeded,ai_audience.refresh.incremental_tick,ai_audience.refresh.daily_tick,ai_audience.member.entered,ai_audience.member.updated,ai_audience.member.exited,ai_audience.run.refreshed
Environment=AICRM_INTERNAL_EVENTS_ALLOWED_EVENT_CONSUMERS=ai_audience.source.changed:ai_audience_source_poke_consumer,channel_entry.entered:ai_audience_source_poke_consumer,questionnaire.submitted:ai_audience_source_poke_consumer,external_form.submitted:ai_audience_source_poke_consumer,payment.succeeded:ai_audience_source_poke_consumer,transaction.paid:ai_audience_source_poke_consumer,payment_succeeded:ai_audience_source_poke_consumer,ai_audience.refresh.incremental_tick:ai_audience_incremental_refresh_consumer,ai_audience.refresh.daily_tick:ai_audience_daily_refresh_consumer,ai_audience.run.refreshed:ai_audience_outbound_effect_planner
Environment=AICRM_AI_AUDIENCE_SCHEDULER_BATCH_SIZE=20
ExecStart=/bin/bash -c 'source /home/ubuntu/venvs/openclaw/bin/activate && python scripts/run_ai_audience_scheduler.py --run-consumers --execute --batch-size "${AICRM_AI_AUDIENCE_SCHEDULER_BATCH_SIZE:-20}"'
