# deploy/openclaw-questionnaire-continuation-scheduler.service
[Unit]
Description=OpenClaw Questionnaire Identity Continuation Reconciliation
After=network.target openclaw-wecom-postgres.service
Requires=openclaw-wecom-postgres.service

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/home/ubuntu/极简 crm
EnvironmentFile=/home/ubuntu/.openclaw-wecom-pg.env
ExecStart=/bin/bash -lc 'source /home/ubuntu/venvs/openclaw/bin/activate && python scripts/run_questionnaire_continuation_scheduler.py --execute --limit "${AICRM_QUESTIONNAIRE_CONTINUATION_BATCH_SIZE:-100}"'
