# deploy/openclaw-broadcast-hourly-feishu-report.service
[Unit]
Description=OpenClaw Broadcast Queue Feishu hourly report
After=network.target openclaw-wecom-postgres.service
Requires=openclaw-wecom-postgres.service

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/home/ubuntu/极简 crm
EnvironmentFile=/home/ubuntu/.openclaw-wecom-pg.env
ExecStart=/bin/bash -lc 'source /home/ubuntu/venvs/openclaw/bin/activate && python scripts/run_broadcast_hourly_feishu_report.py'
