# deploy/openclaw-wecom-callback-ingress.service
[Unit]
Description=OpenClaw WeCom Callback Ingress - fast ACK into webhook inbox
After=network.target openclaw-wecom-postgres.service
Requires=openclaw-wecom-postgres.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/极简 crm
EnvironmentFile=/home/ubuntu/.openclaw-wecom-pg.env
Environment=WECOM_CALLBACK_INGRESS_HOST=127.0.0.1
Environment=WECOM_CALLBACK_INGRESS_PORT=5002
ExecStart=/bin/bash -lc 'source /home/ubuntu/venvs/openclaw/bin/activate && python scripts/run_wecom_callback_ingress.py'
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
