# deploy/openclaw-internal-event-worker.service
[Unit]
Description=OpenClaw Internal Event Queue Worker - relay transactional outbox and dispatch due consumers
After=network.target openclaw-wecom-postgres.service
Requires=openclaw-wecom-postgres.service

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/home/ubuntu/极简 crm
EnvironmentFile=/home/ubuntu/.openclaw-wecom-pg.env
Environment=AICRM_INTERNAL_EVENTS_ENABLED=1
Environment=AICRM_INTERNAL_EVENTS_PAYMENT_ENABLED=1
Environment=AICRM_INTERNAL_EVENTS_QUESTIONNAIRE_ENABLED=1
Environment=AICRM_INTERNAL_EVENTS_SHADOW_ONLY=1
Environment=AICRM_INTERNAL_EVENTS_AUTO_EXECUTE=1
Environment=AICRM_INTERNAL_EVENT_RELAY_ROLE=owner
Environment=AICRM_INTERNAL_EVENT_WORKER_BATCH_SIZE=50
Environment=AICRM_INTERNAL_EVENTS_WORKER_BATCH_SIZE=50
Environment=AICRM_INTERNAL_EVENTS_AUTO_EXECUTE_MAX_BATCH_SIZE=50
Environment=AICRM_INTERNAL_EVENTS_ALLOWED_EVENT_TYPES=payment.succeeded,refund.succeeded,questionnaire.submitted
Environment=AICRM_INTERNAL_EVENTS_ALLOWED_EVENT_CONSUMERS=payment.succeeded:order_projection_consumer,payment.succeeded:service_period_entitlement_consumer,payment.succeeded:product_paid_wecom_tag_consumer,payment.succeeded:webhook_order_paid_consumer,payment.succeeded:customer_business_summary_consumer,payment.succeeded:dnd_policy_consumer,payment.succeeded:ai_assist_notify_consumer,payment.succeeded:ai_audience_source_poke_consumer,refund.succeeded:service_period_refund_consumer,questionnaire.submitted:questionnaire_projection_consumer,questionnaire.submitted:questionnaire_webhook_consumer,questionnaire.submitted:questionnaire_tag_consumer,questionnaire.submitted:automation_questionnaire_consumer,questionnaire.submitted:customer_summary_consumer
ExecStart=/usr/bin/env AICRM_INTERNAL_EVENTS_ENABLED=1 AICRM_INTERNAL_EVENTS_PAYMENT_ENABLED=1 AICRM_INTERNAL_EVENTS_QUESTIONNAIRE_ENABLED=1 AICRM_INTERNAL_EVENTS_SHADOW_ONLY=1 AICRM_INTERNAL_EVENTS_AUTO_EXECUTE=1 AICRM_INTERNAL_EVENT_WORKER_BATCH_SIZE=50 AICRM_INTERNAL_EVENTS_WORKER_BATCH_SIZE=50 AICRM_INTERNAL_EVENTS_AUTO_EXECUTE_MAX_BATCH_SIZE=50 AICRM_INTERNAL_EVENTS_ALLOWED_EVENT_TYPES=payment.succeeded,refund.succeeded,questionnaire.submitted AICRM_INTERNAL_EVENTS_ALLOWED_EVENT_CONSUMERS=payment.succeeded:order_projection_consumer,payment.succeeded:service_period_entitlement_consumer,payment.succeeded:product_paid_wecom_tag_consumer,payment.succeeded:webhook_order_paid_consumer,payment.succeeded:customer_business_summary_consumer,payment.succeeded:dnd_policy_consumer,payment.succeeded:ai_assist_notify_consumer,payment.succeeded:ai_audience_source_poke_consumer,refund.succeeded:service_period_refund_consumer,questionnaire.submitted:questionnaire_projection_consumer,questionnaire.submitted:questionnaire_webhook_consumer,questionnaire.submitted:questionnaire_tag_consumer,questionnaire.submitted:automation_questionnaire_consumer,questionnaire.submitted:customer_summary_consumer /home/ubuntu/venvs/openclaw/bin/python scripts/run_internal_event_worker.py --execute --limit 50
