# deploy/aicrm-next/openclaw-external-contact-full-sync.timer
[Unit]
Description=Run OpenClaw WeCom external contact full reconciliation daily

[Timer]
OnCalendar=*-*-* 03:20:00
Persistent=true
Unit=openclaw-external-contact-full-sync.service

[Install]
WantedBy=timers.target
