# deploy/aicrm-next/openclaw-automation-member-backfill.timer
[Unit]
Description=Run OpenClaw Campaign-ready automation member backfill hourly

[Timer]
OnCalendar=*-*-* *:00:00
Persistent=true
Unit=openclaw-automation-member-backfill.service

[Install]
WantedBy=timers.target
