# deploy/aicrm-next/openclaw-automation-member-backfill.service
[Unit]
Description=OpenClaw Campaign-ready automation member backfill
After=network.target openclaw-wecom-postgres.service
Requires=openclaw-wecom-postgres.service

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/home/ubuntu/极简 crm
EnvironmentFile=/home/ubuntu/.openclaw-wecom-pg.env
ExecStart=/bin/bash -lc 'source /home/ubuntu/venvs/openclaw/bin/activate && python scripts/run_automation_member_backfill.py --limit 5000'
